# spawn-php-fcgi@.socket
[Unit]
Description=PHP FastCGI socket

[Socket]
ListenStream=%t/spawn-php-fcgi/%i.socket

[Install]
WantedBy=sockets.target
